# backend.service.template
[Unit]
Description=Backend API Service
After=network.target

[Service]
Type=simple
User={{USER}}
Group={{GROUP}}
WorkingDirectory={{PROJECT_DIR}}
Environment="PATH={{PROJECT_DIR}}/.venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart={{PROJECT_DIR}}/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port {{PORT}}
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
